# scripts/systemd/prox-rq-backend.service
[Unit]
Description=prox backend

[Service]
WorkingDirectory=/home/prox/prox-server
ExecStart=/home/prox/prox-server/env/bin/rq worker rq_settings
Restart=always
User=prox
Group=prox

[Install]
WantedBy=multi-user.target
